# rc.local
#!/bin/bash
#
# /etc/rc.local: Local multi-user startup script.
#

powertop --auto-tune
ifconfig wlp3s0 down
ifconfig wlp3s0 up

return 0
